# misc/OLD_DOTFILES_BAK/etc/rc.local
#
# /etc/rc.local
#
# Authors:
#    Iyyel <github.com/iyyel>
#

# Fix dell touch keys
/lib/udev/keymap input/event0 /opt/dell_touchkeys_keymap
